# nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f -L 15

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'

